# manifest/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/with-contenv bash

exec "php-fpm${PHP_VERSION}"
