# manifest/etc/s6-overlay/s6-rc.d/php-fpm-tmpl/up
/etc/s6-overlay/script/setup_php-fpm.sh
